# templates/mq-supervisord-prod.conf
[unix_http_server]
file=%(here)s/supervisor.sock ;

[supervisord]
logfile=%(here)s/supervisor.log ;
logfile_maxbytes=50MB ;
logfile_backups=10 ;
loglevel=info ;
pidfile=%(here)s/supervisor.pid ;

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///%(here)s/supervisor.sock ;


[program:live-smtp-realtime]
directory=%(ENV_HERMES_HOME)s ;
command=pipenv run python %(ENV_HERMES_HOME)s/repos/hermes-server/hermes_jobs/mq --agent-type=live-smtp-realtime --agent-limit=0 --agent-parameter=1 ;
numprocs=1 ;
numprocs_start=1 ;
process_name=%(process_num)02d ;
environment=PYTHONPATH=%(ENV_HERMES_HOME)s/repos/hermes-server:%(ENV_HERMES_HOME)s/repos/hermes-client:%(ENV_HERMES_HOME)s/repos/hermes-cpt:%(ENV_HERMES_HOME)s/repos/hermes-superviseur ;
stdout_logfile=%(ENV_HERMES_HOME)s/ops/logs/mq/stdout--smtp-realtime.log ;
stdout_logfile_backups=5 ;
stdout_logfile_maxbytes=50MB ;
stderr_logfile=%(ENV_HERMES_HOME)s/ops/logs/mq/stderr--smtp-realtime.log ;
stderr_logfile_backups=5 ;
stderr_logfile_maxbytes=50MB ;
priority=400 ;


[program:live-smtp-checker]
directory=%(ENV_HERMES_HOME)s ;
command=pipenv run python %(ENV_HERMES_HOME)s/repos/hermes-server/hermes_jobs/mq --agent-type=live-smtp-checker --agent-limit=0 --agent-parameter=1 ;
numprocs=1 ;
numprocs_start=1 ;
process_name=%(process_num)02d ;
environment=PYTHONPATH=%(ENV_HERMES_HOME)s/repos/hermes-server:%(ENV_HERMES_HOME)s/repos/hermes-client:%(ENV_HERMES_HOME)s/repos/hermes-cpt:%(ENV_HERMES_HOME)s/repos/hermes-superviseur ;
stdout_logfile=%(ENV_HERMES_HOME)s/ops/logs/mq/stdout--smtp-checker.log ;
stdout_logfile_backups=5 ;
stdout_logfile_maxbytes=50MB ;
stderr_logfile=%(ENV_HERMES_HOME)s/ops/logs/mq/stderr--smtp-checker.log ;
stderr_logfile_backups=5 ;
stderr_logfile_maxbytes=50MB ;
priority=400 ;


[program:live-smtp]
directory=%(ENV_HERMES_HOME)s ;
command=pipenv run python %(ENV_HERMES_HOME)s/repos/hermes-server/hermes_jobs/mq --agent-type=live-smtp --agent-limit=0 ;
numprocs=2 ;
numprocs_start=1 ;
process_name=%(process_num)02d ;
environment=PYTHONPATH=%(ENV_HERMES_HOME)s/repos/hermes-server:%(ENV_HERMES_HOME)s/repos/hermes-client:%(ENV_HERMES_HOME)s/repos/hermes-cpt:%(ENV_HERMES_HOME)s/repos/hermes-superviseur ;
stdout_logfile=%(ENV_HERMES_HOME)s/ops/logs/mq/stdout--smtp.log ;
stdout_logfile_backups=5 ;
stdout_logfile_maxbytes=50MB ;
stderr_logfile=%(ENV_HERMES_HOME)s/ops/logs/mq/stderr--smtp.log ;
stderr_logfile_backups=5 ;
stderr_logfile_maxbytes=50MB ;
priority=300 ;


[program:live-monitoring-compute]
directory=%(ENV_HERMES_HOME)s ;
command=pipenv run python %(ENV_HERMES_HOME)s/repos/hermes-server/hermes_jobs/mq --agent-type=live-monitoring-compute --agent-limit=0 ;
numprocs=1 ;
numprocs_start=1 ;
process_name=%(process_num)02d ;
environment=PYTHONPATH=%(ENV_HERMES_HOME)s/repos/hermes-server:%(ENV_HERMES_HOME)s/repos/hermes-client:%(ENV_HERMES_HOME)s/repos/hermes-cpt:%(ENV_HERMES_HOME)s/repos/hermes-superviseur ;
stdout_logfile=%(ENV_HERMES_HOME)s/ops/logs/mq/stdout--monitoring-compute.log ;
stdout_logfile_backups=5 ;
stdout_logfile_maxbytes=50MB ;
stderr_logfile=%(ENV_HERMES_HOME)s/ops/logs/mq/stderr--monitoring-compute.log ;
stderr_logfile_backups=5 ;
stderr_logfile_maxbytes=50MB ;
priority=200 ;


[program:live-monitoring-post-processing]
directory=%(ENV_HERMES_HOME)s ;
command=pipenv run python %(ENV_HERMES_HOME)s/repos/hermes-server/hermes_jobs/mq --agent-type=live-monitoring-post-processing --agent-limit=0 ;
numprocs=1 ;
numprocs_start=1 ;
process_name=%(process_num)02d ;
environment=PYTHONPATH=%(ENV_HERMES_HOME)s/repos/hermes-server:%(ENV_HERMES_HOME)s/repos/hermes-client:%(ENV_HERMES_HOME)s/repos/hermes-cpt:%(ENV_HERMES_HOME)s/repos/hermes-superviseur ;
stdout_logfile=%(ENV_HERMES_HOME)s/ops/logs/mq/stdout--monitoring-post-processing.log ;
stdout_logfile_backups=5 ;
stdout_logfile_maxbytes=50MB ;
stderr_logfile=%(ENV_HERMES_HOME)s/ops/logs/mq/stderr--monitoring-post-processing.log ;
stderr_logfile_backups=5 ;
stderr_logfile_maxbytes=50MB ;
priority=200 ;


[program:live-metrics-conso]
directory=%(ENV_HERMES_HOME)s ;
command=pipenv run python %(ENV_HERMES_HOME)s/repos/hermes-server/hermes_jobs/mq --agent-type=live-metrics-conso --agent-limit=0 ;
numprocs=1 ;
numprocs_start=1 ;
process_name=%(process_num)02d ;
environment=PYTHONPATH=%(ENV_HERMES_HOME)s/repos/hermes-server:%(ENV_HERMES_HOME)s/repos/hermes-client:%(ENV_HERMES_HOME)s/repos/hermes-cpt:%(ENV_HERMES_HOME)s/repos/hermes-superviseur ;
stdout_logfile=%(ENV_HERMES_HOME)s/ops/logs/mq/stdout--metrics-conso.log ;
stdout_logfile_backups=5 ;
stdout_logfile_maxbytes=50MB ;
stderr_logfile=%(ENV_HERMES_HOME)s/ops/logs/mq/stderr--metrics-conso.log ;
stderr_logfile_backups=5 ;
stderr_logfile_maxbytes=50MB ;
priority=200 ;


[program:live-metrics-env]
directory=%(ENV_HERMES_HOME)s ;
command=pipenv run python %(ENV_HERMES_HOME)s/repos/hermes-server/hermes_jobs/mq --agent-type=live-metrics-env --agent-limit=0 ;
numprocs=1 ;
numprocs_start=1 ;
process_name=%(process_num)02d ;
environment=PYTHONPATH=%(ENV_HERMES_HOME)s/repos/hermes-server:%(ENV_HERMES_HOME)s/repos/hermes-client:%(ENV_HERMES_HOME)s/repos/hermes-cpt:%(ENV_HERMES_HOME)s/repos/hermes-superviseur ;
stdout_logfile=%(ENV_HERMES_HOME)s/ops/logs/mq/stdout--metrics-env.log ;
stdout_logfile_backups=5 ;
stdout_logfile_maxbytes=50MB ;
stderr_logfile=%(ENV_HERMES_HOME)s/ops/logs/mq/stderr--metrics-env.log ;
stderr_logfile_backups=5 ;
stderr_logfile_maxbytes=50MB ;
priority=200 ;


[program:live-metrics-pcmdi]
directory=%(ENV_HERMES_HOME)s ;
command=pipenv run python %(ENV_HERMES_HOME)s/repos/hermes-server/hermes_jobs/mq --agent-type=live-metrics-pcmdi --agent-limit=0 ;
numprocs=1 ;
numprocs_start=1 ;
process_name=%(process_num)02d ;
environment=PYTHONPATH=%(ENV_HERMES_HOME)s/repos/hermes-server:%(ENV_HERMES_HOME)s/repos/hermes-client:%(ENV_HERMES_HOME)s/repos/hermes-cpt:%(ENV_HERMES_HOME)s/repos/hermes-superviseur ;
stdout_logfile=%(ENV_HERMES_HOME)s/ops/logs/mq/stdout--metrics-pcmdi.log ;
stdout_logfile_backups=5 ;
stdout_logfile_maxbytes=50MB ;
stderr_logfile=%(ENV_HERMES_HOME)s/ops/logs/mq/stderr--metrics-pcmdi.log ;
stderr_logfile_backups=5 ;
stderr_logfile_maxbytes=50MB ;
priority=200 ;


[program:live-superviseur]
directory=%(ENV_HERMES_HOME)s ;
command=pipenv run python %(ENV_HERMES_HOME)s/repos/hermes-server/hermes_jobs/mq --agent-type=live-superviseur --agent-limit=0 ;
numprocs=1 ;
numprocs_start=1 ;
process_name=%(process_num)02d ;
environment=PYTHONPATH=%(ENV_HERMES_HOME)s/repos/hermes-server:%(ENV_HERMES_HOME)s/repos/hermes-client:%(ENV_HERMES_HOME)s/repos/hermes-cpt:%(ENV_HERMES_HOME)s/repos/hermes-superviseur ;
stdout_logfile=%(ENV_HERMES_HOME)s/ops/logs/mq/stdout--superviseur.log ;
stdout_logfile_backups=5 ;
stdout_logfile_maxbytes=50MB ;
stderr_logfile=%(ENV_HERMES_HOME)s/ops/logs/mq/stderr--superviseur.log ;
stderr_logfile_backups=5 ;
stderr_logfile_maxbytes=50MB ;
priority=200 ;


[program:live-fe]
directory=%(ENV_HERMES_HOME)s ;
command=pipenv run python %(ENV_HERMES_HOME)s/repos/hermes-server/hermes_jobs/mq --agent-type=live-fe --agent-limit=0 ;
numprocs=1 ;
numprocs_start=1 ;
process_name=%(process_num)02d ;
environment=PYTHONPATH=%(ENV_HERMES_HOME)s/repos/hermes-server:%(ENV_HERMES_HOME)s/repos/hermes-client:%(ENV_HERMES_HOME)s/repos/hermes-cpt:%(ENV_HERMES_HOME)s/repos/hermes-superviseur ;
stdout_logfile=%(ENV_HERMES_HOME)s/ops/logs/mq/stdout--fe.log ;
stdout_logfile_backups=5 ;
stdout_logfile_maxbytes=50MB ;
stderr_logfile=%(ENV_HERMES_HOME)s/ops/logs/mq/stderr--fe.log ;
stderr_logfile_backups=5 ;
stderr_logfile_maxbytes=50MB ;
priority=100 ;


[program:live-alert]
directory=%(ENV_HERMES_HOME)s ;
command=pipenv run python %(ENV_HERMES_HOME)s/repos/hermes-server/hermes_jobs/mq --agent-type=live-alert --agent-limit=0 ;
numprocs=1 ;
numprocs_start=1 ;
process_name=%(process_num)02d ;
environment=PYTHONPATH=%(ENV_HERMES_HOME)s/repos/hermes-server:%(ENV_HERMES_HOME)s/repos/hermes-client:%(ENV_HERMES_HOME)s/repos/hermes-cpt:%(ENV_HERMES_HOME)s/repos/hermes-superviseur ;
stdout_logfile=%(ENV_HERMES_HOME)s/ops/logs/mq/stdout--alert.log ;
stdout_logfile_backups=5 ;
stdout_logfile_maxbytes=50MB ;
stderr_logfile=%(ENV_HERMES_HOME)s/ops/logs/mq/stderr--alert.log ;
stderr_logfile_backups=5 ;
stderr_logfile_maxbytes=50MB ;
priority=100 ;

